# systemd/btrfs-restic-backup@.timer
[Unit]
Description=Run Restic to backup a BTRFS subvolume
Requires=%N.service

[Timer]
Unit=%N.service
OnCalendar=daily
Persistent=true
